# TerzoAnno/Amministrazione_Sistemi_Reti_e_Cybersecurity/Esami/2025-06-20/large-file-detector/large-file-detector.service
# first and last name: 
# serial number: 
#
# path: ~/.config/systemd/user/large-file-detector.service
#
# command to enable the service: systemctl --user enable large-file-detector.service
# command to start the service: systemctl --user start large-file-detector.service

[Unit]
Description=large file detector service 

[Service]
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=%h/large-file-detector
ExecStart=/usr/bin/python3 app.py --target %h/docs --size 100 --interval 300 --log %h
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target