# TerzoAnno/Amministrazione_Sistemi_Reti_e_Cybersecurity/Esami/2025-07-11/disk-usage-monitor/disk-usage-monitor.service
# nome e cognome:
# matricola:
#
# path: ~/.config/systemd/user/disk-usage-monitor.service

[Unit]
Description=disk usage monitor service

[Service]
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=%h/disk-usage-monitor
ExecStart=/usr/bin/python3 app.py --partition / --threshold 90